# templates/systemd/systemd.service.erb
# nagios-api
#
# By default these are in: /opt/nagios-api

[Unit]
Description="Nagios api"
After=syslog.target network.target
 
[Service]
Type=simple
ExecStart=<%=@nagios_api_start_command%>
ExecStopPost=/bin/sleep 5
Restart=always
 
[Install]
WantedBy=multi-user.target
